# systemd/user/ssh-agent.service.d/auto-add-keys.conf
[Service]
# После запуска агента автоматически добавить все ключи
ExecStartPost=/bin/bash -c 'sleep 1 && for key in ~/.ssh/id_*; do [ -f "$key" ] && [ "${key##*.}" != "pub" ] && ssh-add "$key" 2>/dev/null; done'
